# ALGO/timer/crawler.cron
0 6-12/3 * 12 * root/Dev/ALGO/timer/crawler.sh