# ALGO/timer/token.cron
*/5 * * * * /root/Dev/ALGO/timer/token.sh

